# packages/www-servers/apache-tomcat/files/systemd/tomcat7.service
[Unit]
Description=Tomcat Web Service
After=syslog.target

[Service]
Type=forking
EnvironmentFile=/etc/conf.d/tomcat.conf
ExecStart=/bin/sh /usr/share/tomcat-7/bin/catalina.sh start security
ExecStop=/bin/sh /usr/share/tomcat-7/bin/catalina.sh stop
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target

