# packages/net-misc/freenx/files/systemd/freenx.service
[Unit]
Description=Allow users to get remote X connections via the freenx nxserver
After=local-fs.target
Requires=sshd.service

[Service]
Type=oneshot
ExecStartPre=-/usr/bin/nxserver --cleanup
ExecStart=/usr/bin/nxserver --start
ExecStop=/usr/bin/nxserver --stop
ExecStopPost=-/usr/bin/nxserver --cleanup
StandardOutput=null
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
