# packages/sys-power/apcupsd/files/systemd/apcupsd.service
[Unit]
Description=APC UPS daemon
After=network.target

[Service]
ExecStartPre=/usr/bin/rm -f /etc/apcupsd/powerfail
ExecStart=/sbin/apcupsd -f "/etc/apcupsd/apcupsd.conf" -P "/run/apcupsd.pid"
Type=forking

[Install]
WantedBy=multi-user.target
